# templates/default/datanode.service.erb
[Unit]
Description = DataNode, blocks server for HDFS.
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node.apache_hadoop.hdfs.user %>

PIDFile = <%= node.apache_hadoop.logs_dir %>/hadoop-<%= node.apache_hadoop.hdfs.user %>-datanode.pid
ExecStart = <%= node.apache_hadoop.base_dir %>/sbin/start-dn.sh
ExecStop = <%= node.apache_hadoop.base_dir %>/sbin/stop-dn.sh
ExecReload = <%= node.apache_hadoop.base_dir %>/sbin/restart-dn.sh
LimitNOFILE = 65000
LimitNPROC = 10000

Type = forking

[Install]
WantedBy = multi-user.target
