# templates/default/namenode.service.erb
[Unit]
Description = NameNode, metadata server for HDFS.
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node.apache_hadoop.hdfs.user %>

PIDFile = <%= node[:hadoop][:logs_dir] %>/hadoop-<%= node[:hdfs][:user] %>-namenode.pid
ExecStart = <%= node[:hadoop][:home] %>/sbin/start-nn.sh
ExecStop = <%= node[:hadoop][:home] %>/sbin/stop-nn.sh
ExecReload = <%= node[:hadoop][:home] %>/sbin/restart-nn.sh

Type = forking

[Install]
WantedBy = multi-user.target
