# templates/default/nodemanager.service.erb
[Unit]
Description = NodeManager. The Processing Nodes for YARN.
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node.hadoop.yarn.user %>

PIDFile = <%= node.hadoop.logs_dir %>/yarn-<%= node.hadoop.yarn.user %>-nodemanager.pid
ExecStart = <%= node.hadoop.home %>/sbin/start-nm.sh
ExecStop = <%= node.hadoop.home %>/sbin/stop-nm.sh
ExecReload = <%= node[:hadoop][:home] %>/sbin/restart-nm.sh

Type = forking

[Install]
WantedBy = multi-user.target
