# templates/default/nodemanager.service.erb
[Unit]
Description = NodeManager. The Processing Nodes for YARN.
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node.apache_hadoop.yarn.user %>

PIDFile = <%= node.apache_hadoop.logs_dir %>/yarn-<%= node.apache_hadoop.yarn.user %>-nodemanager.pid
ExecStart = <%= node.apache_hadoop.base_dir %>/sbin/start-nm.sh
ExecStop = <%= node.apache_hadoop.base_dir %>/sbin/stop-nm.sh
ExecReload = <%= node.apache_hadoop.base_dir %>/sbin/restart-nm.sh
LimitNOFILE = 65000
LimitNPROC = 10000

Type = forking

[Install]
WantedBy = multi-user.target
